# docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/var/run/supervisor/supervisord.pid

[program:php-fpm]
command=php-fpm --nodaemonize
autostart=true
autorestart=true
priority=5
stdout_logfile=/var/log/supervisor/php-fpm.log
stderr_logfile=/var/log/supervisor/php-fpm-error.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB

# Laravel Horizon — gerencia os supervisores de queue definidos em
# config/horizon.php (default/uploads, media, face). Single processo
# que faz fork dos workers; reiniciar via `php artisan horizon:terminate`
# (workers terminam graceful após o job atual) ou supervisorctl restart.
[program:horizon]
command=php /var/www/artisan horizon
autostart=true
autorestart=true
priority=20
stopwaitsecs=3600
stopsignal=TERM
stdout_logfile=/var/log/supervisor/horizon.log
stderr_logfile=/var/log/supervisor/horizon-error.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
directory=/var/www

# Scheduler do Laravel — dispara os Schedule::command() em routes/console.php
# (horizon:snapshot, audit:prune, users:cleanup, FetchLiturgyJob, ...).
# schedule:work é o equivalente foreground de um cron entry (* * * * *).
[program:schedule-worker]
command=php /var/www/artisan schedule:work
autostart=true
autorestart=true
priority=15
stdout_logfile=/var/log/supervisor/schedule-worker.log
stderr_logfile=/var/log/supervisor/schedule-worker-error.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
directory=/var/www

# Watches /var/www/.env for in-place edits. When the host writer touches the
# file (host bind-mount → container sees the change), we run config:clear so
# the next request sees the new values without needing a docker compose
# restart. inotify caveat: many editors (VSCode, JetBrains "safe write") do
# rename-replace instead of in-place writes — that fires move_self/delete_self
# on the original inode and inotifywait on the FILE stops listening. We watch
# the parent dir filtered to .env to survive both write modes.
[program:env-watcher]
command=/bin/sh -c "while inotifywait -e close_write,move,create -q /var/www --include '\\.env$' >/dev/null 2>&1; do echo 'env-watcher: .env changed, clearing config cache'; php /var/www/artisan config:clear 2>&1 || true; done"
autostart=true
autorestart=true
startsecs=0
priority=10
stdout_logfile=/var/log/supervisor/env-watcher.log
stderr_logfile=/var/log/supervisor/env-watcher-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
directory=/var/www

[unix_http_server]
file=/var/run/supervisor/supervisor.sock
chmod=0777

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor/supervisor.sock
